# net-p2p/pybitmessage/files/bitmessage
#!/sbin/runscript

KEYS_DAT="/var/lib/bitmessage/.config/PyBitmessage/keys.dat"
MESSAGES_DAT="/var/lib/bitmessage/.config/PyBitmessage/messages.dat"

depend() {
	need net
}

start() {
	ebegin "Starting Bitmessage"

	if [[ ! -e "${KEYS_DAT}" ]] ; then
		cp -a "${KEYS_DAT}.bootstrap" "${KEYS_DAT}"
	fi
  
	if [[ ! -e "${MESSAGES_DAT}" ]] ; then
		cp -a "${MESSAGES_DAT}.bootstrap" "${MESSAGES_DAT}"
	fi
  
	start-stop-daemon --start --user bitmessage:bitmessage --exec /usr/bin/pybitmessage --pidfile /var/run/bitmessage.pid --background --make-pidfile --stdout /var/log/bitmessage/bitmessage.log --stderr /var/log/bitmessage/bitmessage.log
	sleep 10
	eend $?
}

stop() {
	ebegin "Stopping Bitmessage"
	start-stop-daemon --stop --pidfile /var/run/bitmessage.pid
	eend $?
}
